# roles/netkernel/files/autossh.conf
# autossh
 
description "autossh daemon"
 
start on net-device-up IFACE=eth0 
stop on runlevel [01S6] 
 
respawn
respawn limit 5 60 # respawn max 5 times in 60 seconds
 
script
    export AUTOSSH_PIDFILE=/var/run/autossh.pid
    export AUTOSSH_POLL=60
    export AUTOSSH_FIRST_POLL=30
    export AUTOSSH_GATETIME=0
    export AUTOSSH_DEBUG=1
    autossh -M 0 -4 -N -L 10151:localhost:10151 netkernel@license.netkernelroc.org -p 2222 -o "ServerAliveInterval 60" -o "ServerAliveCountMax 3" -o BatchMode=yes -o StrictHostKeyChecking=no -i /home/netkernel/.ssh/id_rsa -2  -vvv
end script
